# roles/miner/templates/miner.timer.j2
[Unit]
Description=Runs miner.service at startup

[Timer]
OnStartupSec=15

[Install]
WantedBy=timers.target
